# ups-oled.service
[Unit]
Description=Python UPS OLED NUT Monitor

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /usr/lib/ups-oled/nut.py
Restart=on-failure
WorkingDirectory=/usr/lib/ups-oled

[Install]
WantedBy=default.target
